# roles/service/monitoring/files/tc_check_update.timer
[Unit]
Description=Run tc_check_update

[Timer]
OnBootSec=30sec
OnUnitActiveSec=6hour

[Install]
WantedBy=timers.target